# systemd/rpcbind.service.in
[Unit]
Description=RPC Bind

[Service]
Type=notify
EnvironmentFile=-/etc/sysconfig/rpcbind
ExecStart=@bindir@/rpcbind $RPCBIND_OPTIONS -w -f

[Install]
WantedBy=multi-user.target
Also=rpcbind.socket
