# systemd/rpcbind.service.in
[Unit]
Description=RPC Bind

[Service]
Type=notify
ExecStart=@bindir@/rpcbind -w -f

[Install]
WantedBy=multi-user.target
Also=rpcbind.socket
